# resources/services/farmconnect-sensor.service
[Unit]
Description=farmconnect-sensor
After=syslog.target network.target redis.service postgresql.service

[Service]
WorkingDirectory=/farmconnect/sensor

ExecStart=/bin/node /farmconnect/sensor/main.js
ExecStop=/bin/kill -15 $MAINPID

RestartSec=10
Restart=always

[Install]
WantedBy=multi-user.target
