# resources/services/farmconnect-sensor-processor.service
[Unit]
Description=farmconnect-sensor-processor
After=syslog.target network.target redis.service

[Service]
WorkingDirectory=/data/farmconnect/sensor-processor

ExecStart=/bin/node /data/farmconnect/sensor-processor/main.js
ExecStop=/bin/kill -15 $MAINPID

RestartSec=10
Restart=always

[Install]
WantedBy=multi-user.target
