# resources/services/farmconnect-ai.service
[Unit]
Description=farmconnect-ai
After=syslog.target network.target redis.service

[Service]
WorkingDirectory=/data/farmconnect/ai

ExecStart=/bin/node /data/farmconnect/ai/main.js
ExecStop=/bin/kill -15 $MAINPID

RestartSec=10
Restart=always

[Install]
WantedBy=multi-user.target
